# roles/netcat-service/templates/netcat.service.j2
[Unit]
Description=Netcat listener on port {{ netcat_port }}
After=network.target
Wants=network.target

[Service]
Type=simple
User={{ netcat_user }}
Group={{ netcat_group }}
ExecStart=/bin/nc -l -k -p {{ netcat_port }}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
